# collector/deploy/tokei-collector.timer.template
[Unit]
Description=Run Tokei Collector every 15 minutes
Requires=tokei-collector.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Unit=tokei-collector.service

[Install]
WantedBy=timers.target
